# dbus/app-org.kde.spectacle.service.in
[Unit]
Description=Spectacle screenshot capture utility

[Service]
ExecStart=@KDE_INSTALL_FULL_BINDIR@/spectacle --dbus
BusName=org.kde.Spectacle
KillMode=process
