# dbus/app-org.kde.spectacle.service.in
[Unit]
Description=Spectacle screenshot capture utility

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/spectacle --dbus
BusName=org.kde.Spectacle
KillMode=process
